# general/home-josh-media-macbook.automount
[Unit]
Description=JP-Ellis Macbook SSH Automount

[Install]
WantedBy=default.target

[Automount]
Where=/home/josh/media/macbook/
TimeoutIdleSec=300
