# general/rustup-update.service
[Unit]
Description=Update Rustup

[Service]
Type=oneshot
ExecStart=/usr/bin/rustup update
