# general/rustup-update.service
[Unit]
Description=Update Rustup
Requisite=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rustup update