# zsh/sheldon-update.service
[Unit]
Description=Update Sheldon plugins

[Service]
Type=oneshot
ExecStart=/usr/bin/sheldon lock --update
